# sys_config/open_route.conf
[program:open_route]
command=/home/sean/open_route/sys_config/gunicorn_start
user=sean
stdout_logfile=/home/sean/open_route/logs/gunicorn_supervisor.log
redirect_stderr=true
