# sys_config/open_route.conf
[program:open_route]
command=~/open_route/sys_config/gunicorn_start
user=ubuntu
stdout_logfile=~/open_route/logs/gunicorn_supervisor.log
redirect_stderr=true
